# modules/civicrm/files/procedures/cron/crontab
*/5 *  *   *   *     sh /var/local/brisskit/brisskit/bin/scheduleTransfer.sh
