# modules/civicrm/files/bk_cron/crontab
* * * * * /bin/bash ~/bk_cron/test_cron.sh >> ~/bk_cron/cron.log
